# prombench/cmd/prombench/supervisord.ini
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
user=root
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:init]
command=/usr/bin/init.sh
exitcodes=0
# logging
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
# no restart
autorestart=false
startretries=0

[program:deploy]
command=/usr/bin/deploy.sh
# logging
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
# no restart
autorestart=false
startretries=0
# no automatic start
autostart=false

[program:clean]
command=/usr/bin/clean.sh
# logging
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
# no restart
autorestart=false
startretries=0
# no automatic start
autostart=false

[eventlistener:signal-supervisor]
command=/usr/bin/signal-supervisor.sh
events=PROCESS_STATE_EXITED,PROCESS_STATE_FATAL
# logging
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0